# docs/systemd/v2raye-server.service
[Unit]
Description=v2rayE Backend API Service
After=network-online.target nss-lookup.target
Wants=network-online.target nss-lookup.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/v2rayE
ExecStart=/usr/lib/v2raye/v2raye --server --api-addr 0.0.0.0:18000 --data-dir /opt/v2rayE
Restart=always
RestartSec=3
Environment=V2RAYN_API_TOKEN=
Environment=V2RAYN_API_ALLOW_PUBLIC=0
Environment=V2RAYN_XRAY_CMD=xray
Environment=V2RAYN_DESKTOP_USER=
Environment=V2RAYN_DATA_DIR=/opt/v2rayE

[Install]
WantedBy=multi-user.target
